# src/bd_track/templates/systemd/bd-track-cleanup.timer.tmpl
[Unit]
Description=Run bd-track Beads cleanup on a regular schedule

[Timer]
OnCalendar=${CHECK_INTERVAL}
Persistent=true
Unit=${SERVICE_NAME}.service

[Install]
WantedBy=timers.target
