# src/bd_track/templates/systemd/bd-track-cleanup.service.tmpl
[Unit]
Description=bd-track Beads cleanup (commit, compact, GC) for all registered projects
After=network.target

[Service]
Type=oneshot
ExecStart=${BD_TRACK_PATH} run-service
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
