# src/bd_timew/templates/systemd/bd-timew-idle-stop.timer.tmpl
[Unit]
Description=Run bd-timew idle server stop check hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=${SERVICE_NAME}.service

[Install]
WantedBy=timers.target
