# src/bd_timew/templates/systemd/bd-timew-cleanup.timer.tmpl
[Unit]
Description=Run bd-timew Beads cleanup on a regular schedule

[Timer]
OnCalendar=${CHECK_INTERVAL}
Persistent=true
Unit=${SERVICE_NAME}.service

[Install]
WantedBy=timers.target
